# debian/postinst
#!/bin/bash -e

echo "Enabling ODAS service"
#systemctl enable matrixio-odas
